# ansible/roles/containers/portainer-agent/templates/portainer-agent-firewall.service.j2
[Unit]
Description=Lock the Portainer agent port to the server in the DOCKER-USER chain
After=docker.service
Requires=docker.service
# PartOf: re-run this oneshot when docker.service restarts (e.g. a docker-ce upgrade), not just on
# reboot — Docker rebuilds DOCKER-USER on restart, so the lock must be re-asserted afterward.
PartOf=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/bash /etc/portainer-agent-firewall.sh

[Install]
WantedBy=multi-user.target
