# ansible/roles/setup/gitops_deploy/templates/gitops-deploy.service.j2
[Unit]
Description=GitOps deployer (pull origin/master, deploy changed services)
After=network-online.target docker.service
Wants=network-online.target
OnFailure=gitops-deploy-alert.service

[Service]
Type=oneshot
# Kill a wedged run (hung ansible/docker) so it can't hold the unit "active" forever and starve
# the 30-min timer (OnUnitActiveSec) — without this a hang stalls the whole pipeline until manual
# intervention. 25min keeps it under the 30-min interval so the next tick still fires. The deployer
# also self-budgets its health-gate phase (RUN_BUDGET_S, default 20min) to leave ~5min here for the
# rollback redeploy — so a slow multi-service gate can't get SIGTERMed before the rollback + hold run.
TimeoutStartSec=25min
User={{ sys_user }}
WorkingDirectory=/home/{{ sys_user }}/server
# Need ~/.local/bin on PATH for `uv` (the deployer runs `uv run ansible-playbook`,
# the same repo-pinned env the operator uses); systemd's default PATH omits it.
Environment=PATH=/home/{{ sys_user }}/.local/bin:/usr/local/bin:/usr/bin:/bin
# flock serializes against the weekly secret-rotate cron (secret-rotate.sh) — both mutate
# this same git working tree, and this deployer's `git reset --hard` on a failed health gate
# could silently discard a half-rotated secret. RUN_START (and so RUN_BUDGET_S, see
# gitops_deploy.py) is measured from process start, AFTER flock acquires — so time spent
# waiting here is pure overhead against this unit's 25min TimeoutStartSec, uncounted by that
# budget. -w 180 outlasts any normal deploy-step contention (secret-rotate is short — 0-2
# tokens) but fails fast enough to leave the health-gate/rollback budget intact; on timeout
# flock exits 1, failing the unit cleanly (OnFailure alerts) and the next 30-min timer tick
# retries rather than racing. /run/lock is world-writable+sticky, so the unprivileged
# {{ sys_user }} can create the lockfile.
ExecStart=/usr/bin/flock -w 180 /var/lock/server-git-tree.lock /usr/bin/python3 /opt/gitops-deploy/gitops_deploy.py
# SOPS age key is already in the user's home; filter plugins auto-load next to
# ansible/deploy.yml, so no ANSIBLE_CONFIG is needed.
