# ansible/roles/setup/gitops_deploy/templates/gitops-deploy.timer.j2
[Unit]
Description=Run the GitOps deployer every 30 minutes

[Timer]
OnBootSec=10min
OnUnitActiveSec=30min
RandomizedDelaySec=2min
Persistent=true

[Install]
WantedBy=timers.target
