# ansible/roles/setup/gitops_deploy/templates/gitops-deploy.service.j2
[Unit]
Description=GitOps deployer (pull origin/master, deploy changed services)
After=network-online.target docker.service
Wants=network-online.target
OnFailure=gitops-deploy-alert.service

[Service]
Type=oneshot
# Kill a wedged run (hung ansible/docker) so it can't hold the unit "active" forever and starve
# the 30-min timer (OnUnitActiveSec) — without this a hang stalls the whole pipeline until manual
# intervention. 25min keeps it under the 30-min interval so the next tick still fires.
TimeoutStartSec=25min
User={{ sys_user }}
WorkingDirectory=/home/{{ sys_user }}/server
# Need ~/.local/bin on PATH for `uv` (the deployer runs `uv run ansible-playbook`,
# the same repo-pinned env the operator uses); systemd's default PATH omits it.
Environment=PATH=/home/{{ sys_user }}/.local/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/usr/bin/python3 /opt/gitops-deploy/gitops_deploy.py
# SOPS age key is already in the user's home; filter plugins auto-load next to
# ansible/deploy.yml, so no ANSIBLE_CONFIG is needed.
