# ansible/roles/setup/renovate_notify/templates/renovate-notify-alert.service.j2
[Unit]
Description=Alert when renovate-notify fails

[Service]
Type=oneshot
# Webhook is EMBEDDED (rendered from SOPS at deploy), NOT sourced from /etc/renovate-notify/config.env:
# config.env is exactly what can be broken/unreadable when the notifier fails, and a hard
# EnvironmentFile= on that same file would leave THIS alert unable to start → no page. curl's own UA
# is fine (Cloudflare only 1010-blocks python-urllib, not curl). This unit is therefore installed mode
# 0600 (see tasks/main.yml) so the embedded webhook isn't world-readable. Mirrors gitops-deploy-alert.service.
ExecStart=/usr/bin/curl -fsS -m 10 -H "Content-Type: application/json" \
  -d '{"content":"🚨 renovate-notify unit failed on daniel-server — check `journalctl -u renovate-notify`."}' \
  "{{ gitops_deploy_discord_webhook }}"
