# ansible/roles/setup/renovate_notify/templates/renovate-notify.timer.j2
[Unit]
Description=Run the Renovate manual-action notifier daily

[Timer]
OnCalendar=*-*-* 13:00:00
RandomizedDelaySec=5min
Persistent=true

[Install]
WantedBy=timers.target
