# ansible/roles/containers/traefik/templates/docker-user-rules.service.j2
[Unit]
Description=Restrict ports 80/443 to Cloudflare IPs in DOCKER-USER chain
After=docker.service
Requires=docker.service
# PartOf: re-run this oneshot when docker.service restarts (e.g. a docker-ce upgrade), not just on
# reboot — Docker rebuilds DOCKER-USER on restart, so the origin lock must be re-asserted afterward.
PartOf=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/bash /etc/docker-user-rules.sh

[Install]
WantedBy=multi-user.target
